# ansible/roles/wikijs/templates/wikijs.service.j2
[Unit]
# Le service démarre après le réseau, car il dépend de PostgreSQL distant.
Description=Wiki.js
After=network.target

[Service]
# Wiki.js est lancé en processus simple depuis le répertoire d'installation.
Type=simple
ExecStart=/usr/bin/node server
WorkingDirectory=/opt/wikijs
Restart=always
Environment=NODE_ENV=production
User=root

[Install]
# L'activation au boot maintient le site disponible après redémarrage de la VM.
WantedBy=multi-user.target
